# packaging/deb/debian/postinst
#!/bin/bash
set -e

# postinst script for p2p-platform

case "$1" in
    configure)
        # Create p2p user if not exists
        if ! getent passwd p2p > /dev/null; then
            adduser --system --group --home /opt/p2p-platform \
                    --no-create-home --disabled-login \
                    --gecos "P2P Platform Service" p2p
        fi

        # Create directories
        mkdir -p /var/log/p2p-platform
        mkdir -p /var/lib/p2p-platform

        # Install Python dependencies
        echo "Installing Python dependencies..."
        cd /opt/p2p-platform/stun-server && pip3 install -r requirements.txt --quiet
        cd /opt/p2p-platform/relay-server && pip3 install -r requirements.txt --quiet
        cd /opt/p2p-platform/signaling-server && pip3 install -r requirements.txt --quiet
        cd /opt/p2p-platform/did-service && pip3 install -r requirements.txt --quiet

        # Set permissions
        chown -R p2p:p2p /opt/p2p-platform
        chown -R p2p:p2p /var/log/p2p-platform
        chown -R p2p:p2p /var/lib/p2p-platform

        # Reload systemd
        systemctl daemon-reload

        # Enable services (but don't start them automatically)
        systemctl enable stun-server.service
        systemctl enable relay-server.service
        systemctl enable signaling-server.service
        systemctl enable did-service.service

        echo "P2P Platform installed successfully!"
        echo "Configure services in /etc/p2p-platform/"
        echo "Start services with: systemctl start <service-name>"
        ;;

    abort-upgrade|abort-remove|abort-deconfigure)
        ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
        ;;
esac

exit 0
