# packaging/deb/debian/prerm
#!/bin/bash
set -e

# prerm script for p2p-platform

case "$1" in
    remove|upgrade|deconfigure)
        # Stop services
        systemctl stop stun-server.service || true
        systemctl stop relay-server.service || true
        systemctl stop signaling-server.service || true
        systemctl stop did-service.service || true

        # Disable services
        systemctl disable stun-server.service || true
        systemctl disable relay-server.service || true
        systemctl disable signaling-server.service || true
        systemctl disable did-service.service || true
        ;;

    failed-upgrade)
        ;;

    *)
        echo "prerm called with unknown argument \`$1'" >&2
        exit 1
        ;;
esac

exit 0
